# runit/sv/runsvdir-david/finish
#!/bin/sh

sv -w600 force-stop /home/david/service/*
sv exit /home/david/service/*
