# dropbox/.config/systemd/user/dropbox.service.d/override.conf
# Dropbox exits with a failure status when the Wayland session goes away at
# logout ("Error reading events from display: Broken pipe"). The systemd --user
# manager outlives the graphical session, so Restart=on-failure -- from the
# packaged /usr/lib/systemd/user/dropbox.service -- resurrects Dropbox within a
# second, during the gap before Hyprland has started waybar.
#
# That matters because Dropbox registers its tray icon via Ayatana/appindicator,
# which requires waybar's org.kde.StatusNotifierWatcher to already own its D-Bus
# name. Dropbox never retries, so an early start means no tray icon all session.
#
# Let it stay dead at logout instead; hyprland.lua starts it after waybar.
[Service]
Restart=no
