# system/samba/etc/systemd/system/smb.service.d/require-media-mount.conf
[Unit]
# Don't start smbd until /mnt/media is mounted, so the [Photos] share is
# never served from an empty, unmounted directory. Adds both After= and
# Requires= on mnt-media.mount.
RequiresMountsFor=/mnt/media
