# Library/LaunchAgents/david.syncmail.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>david.syncmail</string>
        <key>ProgramArguments</key>
        <array>
            <string>sh</string>
            <string>-c</string>
            <string>"${HOME}/bin/syncmail"</string>
        </array>
        <key>StandardOutPath</key>
        <string>/tmp/syncmail.log</string>
        <key>StandardErrorPath</key>
        <string>/tmp/syncmail.log</string>
        <key>RunAtLoad</key>
        <true/>
        <key>StartInterval</key>
        <integer>600</integer>
    </dict>
</plist>

