# notmuch/.config/systemd/user/notmuch.service
[Unit]
Description=Notmuch new mail incorporation
Documentation=man:notmuch-new(1)

[Service]
Type=oneshot
ExecStart=/usr/bin/notmuch new