# restic/.config/systemd/user/backup.service
[Unit]
Description=Backup run
Documentation=man:restic-backup(1)

[Service]
Type=oneshot
ExecStart=%h/bin/backup.sh daily
